# docker/supervisor/elasticsearch.conf
[program:elasticsearch]
directory=/
command=/root/elasticsearch-1.1.1/bin/elasticsearch
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/elasticsearch.out.log
stderr_logfile=/var/log/supervisor/elasticsearch.err.log
